# nginxlog.conf
[program:nginxlog-agent]
user=root
directory=/usr/local/bin
command=/usr/local/bin/nginxlog-agent -f /etc/nginxlog-agent.yaml
autorestart=true
stopwaitsecs=86400
